# .config/systemd/user/backup.service
# https://wiki.archlinux.org/title/Rsync

[Unit]
Description=Backs up files

[Service]
ExecStart=/usr/bin/rsync %h/Sync -CERrltm pi:Sync
